# ws/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=python server.py

directory=/app/server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0